# files/oscam.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=95
STOP=10

SERVICE_NAME=oscam
CONFIG=/etc/config/oscam
BIN=/usr/bin/oscam

start_service() {
    procd_open_instance
    procd_set_param command $BIN -c /etc/oscam
    procd_set_param respawn 300 5
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
}

stop_service() {
    killall oscam
}
